# deploy/kleinanzeige-web.service
[Unit]
Description=Kleinanzeige Web
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kleinanzeige/apps/web
ExecStart=/usr/bin/node ./src/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/opt/kleinanzeige/.env

[Install]
WantedBy=multi-user.target
